# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for myblog"
start on net-device-up
stop on shutdown
respawn
setuid ubuntu
chdir /home/ubuntu/myblog/source
exec ../virtualenv/bin/gunicorn myblog.wsgi:application
